# ansible/roles/setup_node_exporter_raw/templates/node-exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
ExecStart=/opt/prometheus/node-exporter --collector.textfile.directory /opt/prometheus/file_collector

[Install]
WantedBy=multi-user.target
